# debian/postinst
#! /bin/sh
# postinst script for xivo-web-interface
#
# see: dh_installdeb(1)

WEBI_LOG="/var/log/xivo-web-interface"
LIST_DIR_TOCLEAN="/var/lib/pf-xivo-web-interface
					/var/log/pf-xivo-web-interface
					/etc/pf-xivo/web-interface"
LIST_FILE_TOCLEAN="/etc/pf-xivo/web-interface/location.ini
					/etc/pf-xivo/web-interface/report.ini
					/etc/pf-xivo/web-interface/template.ini
					/etc/cron.d/xivo-web-interface
					/etc/default/xivo-agent
					/etc/xivo/web-interface/cti.ini
					/var/lib/xivo/configured"

case "$1" in
	configure)
		chown -R www-data:www-data ${WEBI_LOG}

		old_webi_config_dir="/etc/pf-xivo/web-interface"
		for file in xivo ipbx cti; do
			ini_file="${old_webi_config_dir}/${file}.ini"
			if [ -f ${ini_file}.xivo-upgrade-old ]; then
				mv ${ini_file}.xivo-upgrade-old ${ini_file}
			fi
		done

		for DIR in ${LIST_DIR_TOCLEAN}; do
			if [ -d "${DIR}" ]; then
				rm -rf "${DIR}"
			fi
		done

		for FILE in ${LIST_FILE_TOCLEAN}; do
			if [ -f "${FILE}" ]; then
				rm -f "${FILE}"
			fi
		done

		if dpkg --compare-versions "$2" le '16.13~20161004.141420.044d600'; then
			ln -sf /etc/nginx/locations/http-available/xivo-web-interface \
			       /etc/nginx/locations/http-enabled/xivo-web-interface
			ln -sf /etc/nginx/locations/https-available/xivo-web-interface \
			       /etc/nginx/locations/https-enabled/xivo-web-interface
		fi

		if dpkg --compare-versions "$2" le '18.01'; then
			rm -f /etc/php5/cgi/conf.d/xivo.ini
		fi

		if [ -x /etc/init.d/spawn-fcgi ]; then
			invoke-rc.d spawn-fcgi restart
		fi
	;;
	abort-upgrade|abort-remove|abort-deconfigure)

	;;

	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# must be done after the config file is renamed (automatically generated)
sed -i 's#/var/lib/pf-xivo#/var/lib/xivo#' /etc/xivo/web-interface/ipbx.ini

exit 0
